# miniconda/etc/profile.d/conda.csh
setenv CONDA_EXE "`cygpath 'C:/Users/vjvan/OneDrive/Desktop/GitHub/Text-Extraction-Model/miniconda\Scripts\conda.exe'`";
setenv _CONDA_EXE "`cygpath 'C:/Users/vjvan/OneDrive/Desktop/GitHub/Text-Extraction-Model/miniconda\Scripts\conda.exe'`";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "`cygpath 'C:/Users/vjvan/OneDrive/Desktop/GitHub/Text-Extraction-Model/miniconda\python.exe'`";
setenv _CONDA_ROOT "`cygpath 'C:/Users/vjvan/OneDrive/Desktop/GitHub/Text-Extraction-Model/miniconda'`";

source "`cygpath 'C:/Users/vjvan/OneDrive/Desktop/GitHub/Text-Extraction-Model/miniconda\Lib\site-packages\conda\shell\etc\profile.d\conda.csh'`"